# src/etc/s6-overlay/s6-rc.d/init-permissions/run
#!/command/with-contenv bash
# shellcheck shell=bash

# permissions
chown -R disty:disty \
    /var/run/redis

chown -R disty:disty \
    /config/redis

chown -R disty:disty \
    /etc/redis